# smartctl.rules
 
ACTION=="add|change", KERNEL=="sd[a-z]", RUN+="/bin/smartctl --set=wcache,off /dev/%k"
